# files/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
user = postgres
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
stderr_logfile = /var/log/supervisor/supervisord.log
stdout_logfile = /var/log/supervisor/supervisord.log
autostart=true
autorestart=unexpected
exitcodes=0,2
startsecs=10
